# docker/crontab
*	*	*	* 	* 	/usr/local/bin/set-pin-code.sh 25
